# scripts/generate_events/crontab
0 0 15 * *  /usr/local/bin/node /app/scripts/generate_events/index.js
